# nc-server/supervisord.conf
[supervisord]
logfile = /var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=debug
nodaemon = true
minfds = 1024
minprocs = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nc_server]
process_name	= nc_server
startsecs = 0
directory	= /Debug
command		= sh -c './Shine_srv 5000'
redirect_stderr=true
autorestart	= false

; [program:wazuh-agent]
; process_name	= wazuh-agent
; startsecs = 0
; directory = /
; command		= sh -c './entrypoint.sh'
; redirect_stderr=true
; autorestart	= false